# systemd/gunicorn.service
# Gunicorn systemd service for Cloths by ASF
#
# Install:
#   sudo cp systemd/gunicorn.service /etc/systemd/system/gunicorn-cloths.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable gunicorn-cloths
#   sudo systemctl start gunicorn-cloths
#   sudo systemctl status gunicorn-cloths

[Unit]
Description=Cloths by ASF — Gunicorn WSGI server
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/Cloths-by-ASF/backend
EnvironmentFile=/home/ubuntu/Cloths-by-ASF/backend/.env
Environment=DJANGO_SETTINGS_MODULE=config.settings.production

ExecStart=/home/ubuntu/Cloths-by-ASF/backend/venv/bin/gunicorn \
    config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 120 \
    --access-logfile /var/log/gunicorn/access.log \
    --error-logfile /var/log/gunicorn/error.log

ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
